# scripts/rc.local
# This is an extract from /etc/rc.local, used to start xlReg on the server
# using daemonize

if [ ! -f /tmp/xlReg ]
then
    mkdir -p /tmp/xlReg
    chown ubuntu.ubuntu /tmp/xlreg
fi
/usr/local/sbin/daemonize -a    \
    -c /var/app/xlreg           \
	-l /tmp/xlreg/lockFile	    \
	-p /tmp/xlreg/pidFile	    \
    -u ubuntu                   \
    /home/ubuntu/go/bin/xlReg   \
	$@
